# ostree-native-container/overlay-base/etc/profile.d/systemd-pager.sh
# Tell systemd to not use a pager when printing information
export SYSTEMD_PAGER=cat